# edgedbpkg/edgedb/IDENTIFIER-SLOT.service.in
[Unit]
Description=@@{description}
Documentation=@@{documentation}
After=syslog.target
After=network.target

[Service]
Type=notify

User=edgedb
Group=edgedb

Environment=EDGEDATA=@@{localstatedir}/lib/edgedb/@@{slot}/data/
RuntimeDirectory=%N

ExecStart=@@{bindir}/edgedb-server -D ${EDGEDATA} --runstate-dir=%t/%N
ExecReload=/bin/kill -HUP ${MAINPID}
KillMode=mixed
KillSignal=SIGINT
TimeoutSec=0

[Install]
WantedBy=multi-user.target
